# com.wallbing.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wallbing.weekly</string>

  <key>WorkingDirectory</key>
  <string>/Users/will.ziegler/projects/wallbing-fetcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/will.ziegler/projects/wallbing-fetcher/venv/bin/python3.13</string>
    <string>fetch_weekly.py</string>
  </array>

  <!-- Run every Monday at 9:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Retry on failure (e.g. no network), stop retrying once it succeeds -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- Wait at least 5 minutes between retries -->
  <key>ThrottleInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/tmp/wallbing.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/wallbing.err</string>
</dict>
</plist>
